# hassio-os-agent/hassio-os-agent.initd
#!/sbin/openrc-run

name="Home Assistant OS Agent"
command="/usr/bin/os-agent"
supervisor="supervise-daemon"

depend() {
	need dbus
}
