# hassio-dbus-openrc/hassio-dbus-openrc.initd
#!/sbin/openrc-run

name="Home Assistant DBus bindings"
command="python3"
command_args="-m dbusrc"
supervisor="supervise-daemon"

depend() {
        need dbus
}
